# autostart/blueman.desktop
[Desktop Entry]
Name[oc]=Applet Blueman
Icon=blueman
Exec=blueman-applet
Terminal=false
Type=Application
Categories=
