# autostart/imsettings-start.desktop
[Desktop Entry]
Type=Application
Version=1.0
Comment=Input Method starter
Exec=imsettings-switch -n -q -x
Terminal=false
