# autostart/xfsettingsd.desktop
[Desktop Entry]
Version=1.0
Exec=xfsettingsd
Icon=org.xfce.settings.manager
Terminal=false
Type=Application
StartupNotify=false
#OnlyShowIn=XFCE;
